# bckup.service
[Unit] 
Description=run bckup script 
[Service] 
ExecStart=/$HOME/bin/bckup 
[Install] 
WantedBy=multi-user.target 

